# deployment/aaip-scraper.timer
[Unit]
Description=AAIP Data Scraper Timer - Run every hour
Requires=aaip-scraper.service

[Timer]
# Run every hour at :00 minutes
OnCalendar=hourly
# Run 5 minutes after boot if missed
OnBootSec=5min
# If missed, run immediately
Persistent=true

[Install]
WantedBy=timers.target
