# deployment/aaip-scraper.service
[Unit]
Description=AAIP Data Collection Orchestrator - Hourly data collection job
After=network.target postgresql.service

[Service]
Type=oneshot
User=randy
Group=randy
WorkingDirectory=/home/randy/deploy/aaip-data/scraper

# Use virtual environment
Environment="PATH=/home/randy/deploy/aaip-data/scraper/venv/bin"

# Database connection (update if needed)
Environment="DATABASE_URL=dbname=aaip_data_trend_dev_db"

# Run the orchestrator (runs all data collection scripts)
ExecStart=/home/randy/deploy/aaip-data/scraper/venv/bin/python3 collect_all_data.py

# Timeout for all collectors (15 minutes)
TimeoutStartSec=900

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aaip-scraper

[Install]
WantedBy=multi-user.target
