# deployment/aaip-extended-collectors.service
[Unit]
Description=AAIP Extended Data Collectors (EE, Economy, Labor Market, Job Bank)
After=network.target postgresql.service

[Service]
Type=oneshot
User=randy
WorkingDirectory=/home/randy/deploy/aaip-data/scraper
Environment="PATH=/home/randy/deploy/aaip-data/scraper/venv/bin:/usr/local/bin:/usr/bin"
ExecStart=/home/randy/deploy/aaip-data/scraper/venv/bin/python3 /home/randy/deploy/aaip-data/scraper/collect_extended_data.py

# Timeout and restart settings
TimeoutStartSec=15min
Restart=no

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aaip-extended-collectors

[Install]
WantedBy=multi-user.target
