# deployment/aaip-backend-test.service
[Unit]
Description=AAIP Data Tracker Backend (Test Environment)
After=network.target postgresql.service

[Service]
Type=simple
User=randy
Group=randy
WorkingDirectory=/home/randy/deploy/aaip-data/backend
Environment="PATH=/home/randy/deploy/aaip-data/backend/venv/bin"
Environment="DATABASE_URL=postgresql://localhost/aaip_data_trend_dev_db"
ExecStart=/home/randy/deploy/aaip-data/backend/venv/bin/uvicorn main_enhanced:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aaip-backend-test

[Install]
WantedBy=multi-user.target
